# configs/supervisor.d/cron.ini
[program:cron]
process_name=cron
command = crond -f -l 0 -L /var/log/cron.log
autostart = true
autorestart = true

[eventlistener:cron_exit]
command=supervisor_killer uwsgi
process_name=cron
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL,PROCESS_STATE_STOPPED,PROCESS_STATE_UNKNOWN
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0